# ackermann_bridge_ws/src/ackermann_bridge/config/ackerman_bridge.service
[Unit]
Description=ROS2 Ackermann Bridge Service
After=network-online.target

[Service]
Type=simple
User=${USER}
Environment="ROS_DOMAIN_ID=0"
ExecStart=/opt/ros/humble/bin/ros2 launch ackermann_bridge bridge.launch.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
